# src/main/scripts/etc/systemd/system/webapp.service
[Unit]
Description=Spring Boot Application
After=syslog.target
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/opt/app
ExecStart=/usr/bin/java -jar /opt/app/webapp.war
SuccessExitStatus=143
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
