# src/main/scripts/etc/systemd/system/webapp.service
[Unit]
Description=Spring Boot Application
After=syslog.target network.target
Wants=network-online.target
After=network-online.target

[Service]
User=ec2-user
WorkingDirectory=/opt/app
ExecStart=/usr/bin/java -jar /opt/app/webapp.war
SuccessExitStatus=143
Restart=on-failure
RestartSec=10
TimeoutStartSec=180
TimeoutStopSec=120

# Resource management
LimitNOFILE=65536
Environment="JAVA_OPTS=-Xms512m -Xmx1024m"

[Install]
WantedBy=multi-user.target